# Crontab/test.crontab
*/30 * * * * sh /Users/Chen/Documents/Development/DebugFiles/temp/Crontab/script.sh
